# supervisor/telegram-queue.conf
[program:telegram-queue]
; Queue worker для Telegram уведомлений
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/vverp/artisan queue:work --queue=telegram,notifications --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/vverp/storage/logs/telegram-queue.log
stopwaitsecs=3600
